# supervisor.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi -d --ini uwsgi.ini

[program:nginx]
command=/etc/init.d/nginx restart

